# Home/.config/autostart/firstload.sh.desktop
[Desktop Entry]
Exec=$HOME/System/Scripts/firstload.sh
Icon=application-x-shellscript
Name=firstload.sh
Type=Application
X-KDE-AutostartScript=true
