# HomeDirectory/.config/autostart/memory.sh.desktop
[Desktop Entry]
Exec=/home/bobs/System/Scripts/memory.sh
Icon=application-x-shellscript
Name=memory.sh
Type=Application
X-KDE-AutostartScript=true
